# 2020/tetctf/babyarmshellcode/xinetd.d/babyfmt_xinetd
service babyfmt
{
    disable = no
    type = UNLISTED
    wait = no
    server = /home/babyfmt/run.sh
    socket_type = stream
    protocol = tcp
    user = babyfmt
    only_from =  127.0.0.1
    port = 8888
    flags = REUSE
    per_source = 5
    rlimit_cpu = 3
    nice = 18
}
